# recon/rootfs/etc/services.d/recon/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Recon..."

# Read port from addon config (with fallback to 8765)
PORT=$(bashio::config 'port' 2>/dev/null || echo "8765")
export INGRESS_PORT=${PORT}

# Ingress path from supervisor
INGRESS_PATH=$(bashio::addon.ingress_entry 2>/dev/null || echo "")
export INGRESS_PATH

bashio::log.info "Port: ${PORT}"
bashio::log.info "Ingress path: ${INGRESS_PATH}"

# Try to bring up wireless interfaces
for iface in $(iw dev | grep Interface | awk '{print $2}'); do
    ip link set "$iface" up 2>/dev/null && bashio::log.info "Interface ${iface} up"
done

exec python3 /app/app.py
